# com.henry.skill-insight.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd 任务模板。安装时必须把下面两处绝对路径改成本机实际安装位置：
    1. ProgramArguments 第二项 → <安装目录>/run_skill_insight.sh
    2. StandardOutPath / StandardErrorPath → <安装目录>/skill-log/skill_insight.log
  改完复制到 ~/Library/LaunchAgents/ 再 `launchctl load`。详见 README.md。
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.henry.skill-insight</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>/Users/henry/dev/2_skills/skill-insight/run_skill_insight.sh</string>
	</array>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Weekday</key>
		<integer>1</integer>
		<key>Hour</key>
		<integer>14</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>
	<key>RunAtLoad</key>
	<true/>
	<key>ProcessType</key>
	<string>Background</string>
	<key>LowPriorityBackgroundIO</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/Users/henry/dev/2_skills/skill-insight/skill-log/skill_insight.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/henry/dev/2_skills/skill-insight/skill-log/skill_insight.log</string>
</dict>
</plist>
